# deploy/serverpc/crontab/25_spike_sorting_qc.sh
#!/usr/bin/env bash
set -e
cd ~/Documents/PYTHON/iblscripts/deploy/serverpc/crontab
source ~/Documents/PYTHON/envs/iblenv/bin/activate
python ephys.py _25_spike_sorting_qc $1
